# borg-backup.timer
[Unit]
Description=Daily BorgBackup

[Timer]
OnCalendar=09:00
Persistent=true

[Install]
WantedBy=timers.target
